# scripts/rpi/pi_helper/check_wifi.service
[Unit]
Description=Check and recover Wi-Fi connectivity

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_wifi.sh
